# internet-monitor.service
# /etc/systemd/system/internet-monitor.service
[Unit]
Description=Internet Connection Monitor
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Documents/github/network-monitoring-tools
ExecStart=/usr/bin/node /home/pi/Documents/github/network-monitoring-tools/internet-monitor.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target